# immich-postgres.container
[Unit]
Description=Immich Postgres

[Container]
#Image=docker.io/tensorchord/pgvecto-rs:pg14-v0.2.0@sha256:739cdd626151ff1f796dc95a6591b55a714f341c737e27f045019ceabf8e8c52
Image=docker.io/tensorchord/pgvecto-rs:pg14-v0.2.0
AutoUpdate=registry
Volume=immich-postgres-data.volume:/var/lib/postgresql/data
EnvironmentFile=/etc/containers/systemd/env/immich.env
Exec=postgres -c shared_preload_libraries=vectors.so -c 'search_path="$$user", public, vectors' -c logging_collector=on -c max_wal_size=2GB -c shared_buffers=512MB -c wal_compression=on
Network=immich.network

[Service]
TimeoutStartSec=900
Restart=always

[Install]
# Start by default on boot
WantedBy=multi-user.target default.target
